# rc/zsh/profile.d/golang.zsh
export GOPATH=$HOME/.ghq
_glob_add_path_env PATH $GOPATH/bin

